# root/etc/uci-defaults/30_luci-theme-noodles
#!/bin/sh
uci batch <<-EOF
	set luci.themes.noodles=/luci-static/noodles
	set luci.main.mediaurlbase=/luci-static/noodles
	commit luci
EOF
exit 0
